# services/ipython-notebook/etc/service/ipython-notebook/run
#!/usr/bin/env bash
exec ipython notebook
